# infrastructure/systemd/noytrix-telegram-bot.service
[Unit]
Description=Noytrix Telegram ScamShield Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/noytrix-telegram-bot
ExecStart=/root/noytrix-telegram-bot/venv/bin/python /root/noytrix-telegram-bot/bot.py
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
