# infrastructure/systemd/noytrix-lead-monitor.service
[Unit]
Description=Noytrix Lead Monitor scan

[Service]
Type=oneshot
WorkingDirectory=/root/noytrix-lead-monitor
ExecStart=/root/noytrix-lead-monitor/venv/bin/python /root/noytrix-lead-monitor/main.py
